# launchd/com.tully.aihot-feishu-monitor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tully.aihot-feishu-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-u</string>
    <string>/Users/tully/01work/aihot-feishu-monitor/aihot_feishu_monitor.py</string>
    <string>--env-file</string>
    <string>/Users/tully/01work/aihot-feishu-monitor/.env</string>
    <string>--loop</string>
  </array>

  <!-- 强制直连：监控只访问公网数据源与飞书 webhook，不需要任何代理。
       显式清空代理变量可避免进程启动时继承一个本地代理(如 Clash 127.0.0.1:7897)，
       否则代理一关就会每轮 Connection refused 却因 KeepAlive 不重启而静默失败。 -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>HTTP_PROXY</key><string></string>
    <key>HTTPS_PROXY</key><string></string>
    <key>ALL_PROXY</key><string></string>
    <key>http_proxy</key><string></string>
    <key>https_proxy</key><string></string>
    <key>all_proxy</key><string></string>
    <key>no_proxy</key><string>*</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/tully/01work/aihot-feishu-monitor</string>

  <key>StandardOutPath</key>
  <string>/Users/tully/.aihot-feishu-monitor/monitor.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/tully/.aihot-feishu-monitor/monitor.err.log</string>
</dict>
</plist>
